# templates/docker-system-prune.service.erb
[Unit]
Description=clean up unused images and stopped docker containers
Documentation=man:docker-system-prune(1)
Documentation=man:docker(1)
Requires=docker.service

[Service]
Type=oneshot
# --force does not prompt for confirmation
ExecStart=/usr/bin/docker system prune --force

[Install]
WantedBy=network.target
