# templates/docker-system-prune.timer.erb
[Unit]
Description=Cleanup unused images and stopped docker containers daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
